# systemd/fluxgate-web.service.d/resource-limits.conf
# =============================================================================
# BRCloud FluxGate - systemd Resource Control (service web)
# =============================================================================
# Override systemd pour plafonner CPU/RAM/FD du service web.
# Empeche un service sous attaque d'absorber toutes les ressources.
#
# Installation :
#   mkdir -p /etc/systemd/system/<service>.service.d/
#   cp resource-limits.conf /etc/systemd/system/<service>.service.d/
#   systemctl daemon-reload && systemctl restart <service>
#
# Remplacer <service> par le nom de votre service (nginx, apache2, etc.)
# =============================================================================

[Service]
# --- CPU ---
# Limiter a 80% d'un coeur CPU (ou 200% pour 2 coeurs, etc.)
CPUQuota=80%

# Poids CPU relatif (par rapport aux autres services, defaut 100)
CPUWeight=80

# --- Memoire ---
# Limite dure : le service est OOM-killed s'il depasse
MemoryMax=2G

# Limite haute : pression memoire (reclamation aggressive au-dessus)
MemoryHigh=1800M

# --- File Descriptors ---
# Nombre max de fichiers ouverts (sockets incluses)
LimitNOFILE=65536

# --- Processus ---
# Nombre max de threads/processus
LimitNPROC=4096

# --- Taille max des fichiers core dump (0 = desactive) ---
LimitCORE=0

# --- Securite supplementaire (optionnel) ---
# Empecher le service d'acquerir de nouveaux privileges
NoNewPrivileges=true

# Proteger le systeme de fichiers
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Restreindre les namespaces (si le service n'a pas besoin de containers)
RestrictNamespaces=true

# Restart automatique en cas de crash
Restart=on-failure
RestartSec=5s
